# ros/melodic/x64/etc/catkin/profile.d/10.rosbuild.bat
REM generated from rosbuild/env-hooks/10.rosbuild.bat.em

REM env variables in installspace
if "%CATKIN_ENV_HOOK_WORKSPACE%"=="" set CATKIN_ENV_HOOK_WORKSPACE=C:/opt/ros/melodic/x64
set ROS_ROOT=%CATKIN_ENV_HOOK_WORKSPACE%/share/ros
