# ros/melodic/x64/etc/catkin/profile.d/python_orocos_kdl_site_packages.tcsh
export PYTHONPATH=C:/opt/ros/melodic/x64/Lib\site-packages:$PYTHONPATH
